# roles/system/files/rc.local
#!/usr/bin/env bash
rfkill block bluetooth
exit 0
